# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid


[program:bind]
command=/usr/sbin/named -c /etc/bind/named.conf -f
autostart=true
autorestart=true
